# k3s-master/rootfs/etc/services.d/k3s/finish
#!/usr/bin/with-contenv bashio
bashio::log.info "K3s service stopped."
exit 0
