# k3s-master/rootfs/etc/services.d/k3s/run
#!/usr/bin/with-contenv bashio
set -e

CLUSTER_TOKEN=$(bashio::config 'cluster_token')
DATA_PATH=$(bashio::config 'data_path')
EXTRA_ARGS=$(bashio::config 'extra_args')

mkdir -p "$DATA_PATH"

bashio::log.info "Starting K3s server..."
exec /usr/local/bin/k3s server \
  --token "$CLUSTER_TOKEN" \
  --data-dir "$DATA_PATH" \
  --cluster-init \
  --embedded-registry \
  --disable traefik \
  --disable servicelb \
  --disable local-storage \
  --write-kubeconfig "$DATA_PATH/k3s.yaml" \
  --write-kubeconfig-mode 644 \
  --tls-san "$(hostname -I | awk '{print $1}')" \
  $EXTRA_ARGS
